# bash/crontab
40 12 * * mon-fri sh ~/github/ESP8266/bash/lights-on.sh
